# cron/backup.cron
0 2 * * * /home/ubuntu/scripts/backup-database.sh >> /home/ubuntu/logs/backup.log 2>&1
